# RH124/First_extra_exercice/limpia-secretos.service
[Unit]
Description=RH124 extra homework secrets

[Service]
Type=simple
ExecStart=/usr/bin/limpia-secretos.sh

[Install]
WantedBy=multi-user.target
